# docs/systemd/govcomms.service
[Unit]
Description=GovComms Multi-Bot Service
Documentation=https://github.com/stake-plus/govcomms
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=govcomms
Group=govcomms
WorkingDirectory=/opt/govcomms

# Core configuration
Environment="MYSQL_DSN=govcomms:password@tcp(127.0.0.1:3306)/govcomms"
Environment="DISCORD_TOKEN=your-discord-token"
Environment="GUILD_ID=your-guild-id"
Environment="QA_ROLE_ID=role-id-for-qa"
Environment="RESEARCH_ROLE_ID=role-id-for-research"
Environment="FEEDBACK_ROLE_ID=role-id-for-feedback"
Environment="OPENAI_API_KEY=sk-your-openai-key"
Environment="ENABLE_FEEDBACK=1"

# Optionally load overrides from an env file
EnvironmentFile=-/opt/govcomms/.env.govcomms

ExecStart=/opt/govcomms/bin/govcomms --enable-feedback

Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/govcomms/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

StandardOutput=journal
StandardError=journal
SyslogIdentifier=govcomms

[Install]
WantedBy=multi-user.target
