# docs/systemd/govcomms.service
[Unit]
Description=Stake Plus GovComms Service
Documentation=https://github.com/stake-plus/govcomms
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=govcomms
Group=govcomms
WorkingDirectory=/opt/govcomms
EnvironmentFile=/opt/govcomms/.env.govcomms

# Core execution
ExecStart=/opt/govcomms/bin/govcomms --enable-qa --enable-research --enable-feedback
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/govcomms /var/log/govcomms

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

StandardOutput=journal
StandardError=journal
SyslogIdentifier=govcomms

[Install]
WantedBy=multi-user.target

