# scripts/govcomms-api.service
[Unit]
Description=GovComms API Server
After=network.target mysql.service redis.service

[Service]
Type=simple
User=govcomms
Group=govcomms
WorkingDirectory=/opt/govcomms
ExecStart=/opt/govcomms/bin/api
Restart=always
RestartSec=10
Environment="MYSQL_DSN=govcomms:password@tcp(localhost:3306)/govcomms"
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="PORT=443"
Environment="SSL_CERT=/etc/letsencrypt/live/reeeeeeeeee.io/fullchain.pem"
Environment="SSL_KEY=/etc/letsencrypt/live/reeeeeeeeee.io/privkey.pem"

# Allow access to Let's Encrypt certificates
SupplementaryGroups=ssl-cert

# Logging
StandardOutput=append:/var/log/govcomms/api.log
StandardError=append:/var/log/govcomms/api-error.log

[Install]
WantedBy=multi-user.target